# Zabbix Report App/requirements, service and conf/zabbix-report.service
[Unit]
Description=Zabbix Report App (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
User=root
Group=root
WorkingDirectory=/opt/zabbix_report

Environment="ZABBIX_API_TOKEN=/etc/zabbix_report/token"
Environment="ZABBIX_URL=https://your.domain.name/api_jsonrpc.php"
Environment="PYTHONUNBUFFERED=1"

ExecStart=/usr/bin/python3 -m uvicorn backend.main:app \
    --host 127.0.0.1 \
    --port 8081 \
    --workers 4 \
    --timeout-keep-alive 75

PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ProtectKernelTunables=true
ProtectKernelModules=true
NoNewPrivileges=true
ReadOnlyPaths=/opt/zabbix_report
ReadWritePaths=/opt/zabbix_report/data
ReadWritePaths=/var/log/zabbix-report
RuntimeDirectory=zabbix-report

[Install]
WantedBy=multi-user.target
